# scripts/splat-service.service
[Unit]
Description=Robot Dog 3D Splat Service
After=network.target

[Service]
User=asus
WorkingDirectory=/home/asus/robot-dog-3d-analysis
EnvironmentFile=/home/asus/robot-dog-3d-analysis/.env
ExecStart=/home/asus/robot-dog-3d-analysis/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8090 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
